# event-management/baguette-client/bin/baguette-client
#! /bin/sh

### BEGIN INIT INFO
# Provides:          baguette-client
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: baguette-client
# Description:       Controls the Baguette Client service
### END INIT INFO

export JAVA_HOME="/usr/bin/java"
SU_USER=root
 
startcmd='/opt/baguette-client/bin/run.sh &>>/opt/baguette-client/logs/output.txt &'
stopcmd='/opt/baguette-client/bin/kill.sh &>>/opt/baguette-client/logs/output.txt'
 
case "$1" in
start)
   echo "Starting Baguette Client..."
   su -c "${startcmd}" $SU_USER
;;
restart)
   echo "Re-starting Baguette Client..."
   su -c "${stopcmd}" $SU_USER
   su -c "${startcmd}" $SU_USER
;;
stop)
   echo "Stopping Baguette Client..."
   su -c "${stopcmd}" $SU_USER
;;
*)
   echo "Usage: $0 {start|stop|restart}"
exit 1
esac
